# field-service-bot/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:resident-bot]
command=python3 /app/bot.py
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/null
stderr_logfile=/dev/null

[program:agent-bot]
command=python3 /app/agent_bot.py
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/null
stderr_logfile=/dev/null
